# autostart/aura-welcome.desktop
[Desktop Entry]
Type=Application
Name=Aura Welcome
Comment=Display AuraOS welcome message on first login
Exec=/bin/bash -c 'if [ -f ~/.local/state/aura/first-run.mode ]; then ~/.local/share/aura/install/first-run/welcome.sh && rm ~/.local/state/aura/first-run.mode; fi'
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
